# etc/supervisor/conf.d/fluxbox.conf
[program:fluxbox]
command=/usr/bin/fluxbox -display %(ENV_DISPLAY)s
priority=5
redirect_stderr=true
stdout_logfile_backups=2
stdout_logfile_maxbytes=10MB
