# etc/supervisor/conf.d/tinywm.conf
[program:tinywm]
command=/usr/bin/tinywm
priority=5
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stdout_logfile_backups=4
stdout_logfile_maxbytes=10MB
stderr_logfile=NONE
