# etc/supervisor/conf.d/xvfb.conf
[program:xvfb]
command=/usr/bin/Xvfb %(ENV_DISPLAY)s -screen 0 %(ENV_SCREEN_GEOMETRY)s -ac
priority=0
redirect_stderr=true
stdout_logfile_backups=2
stdout_logfile_maxbytes=10MB
